# crontab
*/5 * * * * bash /path/to/mount/dir/script.sh
